# auto-update.cron
# Update CVE DB every 2 hours, at 33 minutes
33 */2 * * * (/bin/date ; echo "Auto-updating DB" ; /opt/cve/sbin/db_updater.py -v ; /bin/date ; echo "Auto-updating CPE browsing cache" ; /opt/cve/sbin/db_cpe_browser.py ; /bin/date) >> /opt/cve/log/auto-update.log 2>&1
