# auto-update.cron
# Update CVE DB every 2 hours, at 33 minutes
33 */2 * * * /opt/cve/sbin/db_updater.py -v >> /opt/cve/log/auto-update.log 2>&1
